# files/etc/profile.d/baryte-env.sh
# Gaming environment variables
export STEAM_RUNTIME=1
export STEAM_RUNTIME_HEAVY=1

# Development environment setup
export GOPATH=$HOME/go
export PATH=$PATH:$GOPATH/bin